# crontab
*/2 * * * * python3 /tado_data_collection/main.py >> /cron_logfile.log 2>&1
