# deprecated_files/crontab
*/2 * * * * stdbuf -o0 -e0 python3 -u /tado_data_collection/main.py > /dev/stdout
